# deploy/goapi.service
[Unit]
Description=GoApi - Address lookup by CEP (ViaCEP/BrasilAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/GoApi
ExecStart=/root/GoApi/cmd/api/main --server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
